# other/websockify-telnet.socket
[Unit]
Description=Proxy for telnet via websocket

[Socket]
ListenStream=23
Accept=yes

[Install]
WantedBy=sockets.target
